# port-files/pkg-plist
bin/PCDM-session
etc/pcdm.conf.dist
etc/rc.d/pcdm
sbin/PCDMd
share/PCDM/pcdm-session
share/PCDM/themes/default/background.jpg
share/PCDM/themes/default/banner.png
share/PCDM/themes/default/default-desktop.png
share/PCDM/themes/default/default-user.png
share/PCDM/themes/default/default.theme
share/PCDM/themes/default/keyboardlayout.png
share/PCDM/themes/default/localechange.png
share/PCDM/themes/default/login-stealth.png
share/PCDM/themes/default/login.png
share/PCDM/themes/default/nextDE.png
share/PCDM/themes/default/password.png
share/PCDM/themes/default/previousDE.png
share/PCDM/themes/default/splashscreen.png
share/PCDM/themes/default/system.png
share/PCDM/themes/default/virtualkeyboard.png
share/PCDM/themes/enterprise/background.jpg
share/PCDM/themes/enterprise/banner.png
share/PCDM/themes/enterprise/default-desktop.png
share/PCDM/themes/enterprise/default-user.png
share/PCDM/themes/enterprise/enterprise.theme
share/PCDM/themes/enterprise/keyboardlayout.png
share/PCDM/themes/enterprise/localechange.png
share/PCDM/themes/enterprise/login-stealth.png
share/PCDM/themes/enterprise/login.png
share/PCDM/themes/enterprise/nextDE.png
share/PCDM/themes/enterprise/password.png
share/PCDM/themes/enterprise/previousDE.png
share/PCDM/themes/enterprise/splashscreen.png
share/PCDM/themes/enterprise/system.png
share/PCDM/themes/enterprise/virtualkeyboard.png
share/xsessions/_deforaos.desktop
share/xsessions/_evilwm.desktop
share/xsessions/_fluxbox.desktop
share/xsessions/_fvwm-crystal.desktop
share/xsessions/_fvwm.desktop
share/xsessions/_gnome-classic.desktop
share/xsessions/_icewm.desktop
share/xsessions/_kde4.desktop
share/xsessions/_plex.desktop
share/xsessions/_ratpoison.desktop
share/xsessions/_spectrwm.desktop
share/xsessions/_vtwm.desktop
share/xsessions/_windowlab.desktop
share/xsessions/_windowmaker.desktop
