# port-files/files/pcdm-openrc.in
#!/sbin/openrc-run
# Copyright (c) 2016 Joe Maloney <jmaloney@ixsystems.com>
# Released under the 2-clause BSD license.

command=/usr/local/bin/startx
pid=/var/run/PCDMd.pid
command_args="$PCDMd_args /usr/local/bin/PCDM-session"
name="TrueOS Display Manager"

depend() {
	need urandom dbus syscons moused
	after bootmisc
}

start() {
  # Check if boot-loader set option to skip xorg
  if [ "`kenv noxorg 2>/dev/null`" = "YES" ] ; then
     return 0
  fi

  echo "Starting PCDM."

  (
    iter=0
    while [ -e /var/run/nologin ]; do 
        if [ ${iter} -eq 60 ]; then
            break
        fi
        sleep 1
        iter=$(expr ${iter} + 1)
    done

    #Start the PCDM login daemon
    if [ ! -e /var/run/nologin ]; then
      if [ -e /var/tmp/.PCDMstop ]; then 
        rm /var/tmp/.PCDMstop
      fi
      ${command} ${command_args} 
    fi
  ) &
}

stop_post() {
	${command} stop
}
