# port-files/files/pkg-install.in
#!/bin/sh
#
# $FreeBSD$
#

PREFIX=${PKG_PREFIX-/usr/local}

case $2 in
POST-INSTALL)
	if [ -e /sbin/rc-update }; then
	/sbin/rc-update add pcdm default
	fi
	;;
esac

